# deploy/systemd/tredific-api.service
[Unit]
Description=Tredific Django API (Gunicorn)
After=network.target postgresql.service

[Service]
User=tredific
Group=www-data
WorkingDirectory=/home/tredific/tredific/backend
EnvironmentFile=/home/tredific/tredific/backend/.env
ExecStart=/home/tredific/tredific/backend/.venv/bin/gunicorn \
          --workers 3 --threads 2 --timeout 60 \
          --access-logfile /home/tredific/tredific/backend/logs/gunicorn-access.log \
          --error-logfile /home/tredific/tredific/backend/logs/gunicorn-error.log \
          --bind 127.0.0.1:8000 tredific.wsgi:application
Restart=always
RestartSec=3

# Hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/tredific/tredific/backend/media /home/tredific/tredific/backend/logs /home/tredific/tredific/backend/staticfiles
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
